# deploy/sistema-etiquetas-backend.service
[Unit]
Description=Sistema Etiquetas API (Node/Express)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/etiquetas/backend
Environment=NODE_ENV=production
Environment=PORT=50000
Environment=HOST=127.0.0.1
# Descomenta cuando tengas .env en el VPS (SYS21, SQLite, etc.):
# EnvironmentFile=-/var/www/etiquetas/backend/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# Hardening mínimo (ajusta User si tu despliegue usa otro usuario)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
